# templates/default/postgresql.service.erb
[Service]
<% if node['postgresql']['use_pgdg_packages'] %>
.include /usr/lib/systemd/system/postgresql-<%= node['postgresql']['version'] %>.service
<% else %>
.include /usr/lib/systemd/system/postgresql.service
<% end %>

Environment=
Environment=PGPORT=<%= @port %>
Environment=PGDATA=<%= @data_dir %>
